# deploy/systemd/webtaku.service
[Unit]
Description=Webtaku Spring Boot Application
After=network.target mysql.service redis.service

[Service]
Type=simple
User=webtaku
Group=webtaku
WorkingDirectory=/opt/webtaku
EnvironmentFile=/etc/webtaku/webtaku.env
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/webtaku/app.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webtaku

[Install]
WantedBy=multi-user.target
